# jekyll-hook.conf
description "Jekyll-Hook"
author      "Development Seed"

# When to start / stop:
start on filesystem and net-device-up
stop on shutdown

setuid deploy

# Automatically Respawn:
respawn

# What to do:
script
    export HOME="/home/deploy"
    export PATH=$HOME/.nodenv/shims:$HOME/.rbenv/shims:$PATH
    cd $HOME/jekyll-hook
    exec node jekyll-hook.js 2>> $HOME/log/jekyll-hook.log
end script

post-start script
   # Optionally put a script here that will notify you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script

# Source: http://kevin.vanzonneveld.net/techblog/article/run_nodejs_as_a_service_on_ubuntu_karmic/
